# hms-core/etc/init/testdaemon.conf
description "The testdaemon service, a simple service example written in Python"

# Start in normal runlevels when disks are mounted and networking 
start on runlevel [2345]

# Stop on shutdown/halt, single-user mode, reboot
stop on runlevel [016]

env statedir=/var/cache/testdaemon

# Create the cache directory
pre-start /usr/bin/mkdir -p "$statedir"

# Start
exec /home/ballbee/py_envs/hms/core/bin/python /home/ballbee/git/hms_sys/hms-core/src/hms_core/daemons.py
